# trading_bot.service
[Unit]
Description=Money Maker ML Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__LINUX_USER__
WorkingDirectory=__WORKDIR__
ExecStart=__PYTHON_BIN__ main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

[Install]
WantedBy=multi-user.target
